# supervisor/confd.conf
[program:confd]
command=/confd -node="http://${ETCD_HOST}" -verbose=true
priority=10
numprocs=1
autostart=true
autorestart=true