# supervisor/confd.conf
[program:twemproxy]
command=/confd -node="http://${ETCD_HOST}:4001" -quiet=false
priority=10
numprocs=1
autostart=true
autorestart=true