# supervisor/twemproxy.conf
[program:twemproxy]
command=/twemproxy -c /twemproxy.yaml -d
numprocs=1
autostart=true
autorestart=true