# supervisor/twemproxy.conf
[program:twemproxy]
command=/twemproxy -c /twemproxy.yaml -v 8
numprocs=1
autostart=true
autorestart=true
